# kloxo/httpdocs/htmllib/filecore/init.crontab
*   *   *   *   *   sh /usr/local/lxlabs/<%program%>/bin/fetchmailticket.sh
